# imagebuilder/test/qemu-dev/qemu-dev-overlay/etc/uci-defaults/99-qemu-dev
#!/bin/sh
# QEMU dev sandbox — трёхпортовая конфигурация.
# eth0 = WAN (часть br0, DHCP от Docker dnsmasq)
# eth1 = LAN (часть br0)
# eth2 = Management (DHCP, SSH)
# br0 = eth0 + eth1 = прозрачный L2 мост (как в бою)
#
# Прошивка настраивает br0 через setup-bridge.sh.
# eth2 не участвует в bridge — отдельный management.

# BOX_ID для dev sandbox
echo "BB-QEMU-DEV" > /etc/bridgebox/box-id

# Management interface — eth2 (отдельный порт, не в bridge)
echo "eth2" > /etc/bridgebox/mgmt-iface

# Настроить eth2 как management с DHCP
uci set network.mgmt=interface
uci set network.mgmt.device='eth2'
uci set network.mgmt.proto='dhcp'
uci commit network

# SSH на всех интерфейсах
uci set dropbear.@dropbear[0].Interface=''
uci commit dropbear

# Firewall — разрешить всё (dev sandbox)
if uci show firewall >/dev/null 2>&1; then
    uci set firewall.@defaults[0].input='ACCEPT'
    uci set firewall.@defaults[0].output='ACCEPT'
    uci set firewall.@defaults[0].forward='ACCEPT'
    uci commit firewall
fi

exit 0
