# imagebuilder/files/etc/uci-defaults/10-bridgebox-system
#!/bin/sh
# 10-bridgebox-system — Базовая настройка системы при первой загрузке
# Выполняется один раз (uci-defaults).

BRIDGEWRT_VERSION="0.1.0"

# --- Общие утилиты ---
. /usr/lib/bridgebox/lib-common.sh

# --- BOX_ID (генерация из MAC eth0 при первой загрузке) ---
BOX_ID_FILE="/etc/bridgebox/box-id"
EXISTING_ID=$(cat "$BOX_ID_FILE" 2>/dev/null)
if [ -z "$EXISTING_ID" ] || [ "$EXISTING_ID" = "TEMPLATE" ]; then
    # Ждём eth0 (до 10 сек — для rockchip/R3S драйвер грузится асинхронно)
    MAC=""
    attempts=0
    while [ $attempts -lt 10 ]; do
        MAC=$(cat /sys/class/net/eth0/address 2>/dev/null)
        [ -n "$MAC" ] && [ "$MAC" != "00:00:00:00:00:00" ] && break
        MAC=""
        attempts=$((attempts + 1))
        sleep 1
    done
    if [ -n "$MAC" ]; then
        BOX_ID="BB-$(echo -n "$MAC" | sha256sum | cut -c1-6 | tr 'a-f' 'A-F')"
        mkdir -p /etc/bridgebox
        safe_write "$BOX_ID_FILE" "$BOX_ID"
    fi
fi

# --- Branding ---
cat > /etc/banner <<'BANNER'
░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░
░       ░░░       ░░░        ░░       ░░░░      ░░░        ░░  ░░░░  ░░       ░░░        ░
▒  ▒▒▒▒  ▒▒  ▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒  ▒▒  ▒▒▒▒▒▒▒▒  ▒▒▒▒▒▒▒▒  ▒  ▒  ▒▒  ▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒
▓       ▓▓▓       ▓▓▓▓▓▓  ▓▓▓▓▓  ▓▓▓▓  ▓▓  ▓▓▓   ▓▓      ▓▓▓▓        ▓▓       ▓▓▓▓▓▓  ▓▓▓▓
█  ████  ██  ███  ██████  █████  ████  ██  ████  ██  ████████   ██   ██  ███  ██████  ████
█       ███  ████  ██        ██       ████      ███        ██  ████  ██  ████  █████  ████
██████████████████████████████████████████████████████████████████████████████████████████

 ---------------------------------------------------------
 Transparent L2 Bridge + Remote Management
 ---------------------------------------------------------
BANNER

BOX_ID=$(cat /etc/bridgebox/box-id 2>/dev/null)
if [ -n "$BOX_ID" ] && [ "$BOX_ID" != "TEMPLATE" ]; then
    echo " Box: $BOX_ID  |  v${BRIDGEWRT_VERSION}" >> /etc/banner
else
    echo " v${BRIDGEWRT_VERSION}" >> /etc/banner
fi
echo "" >> /etc/banner

sed -i "s/DISTRIB_ID=.*/DISTRIB_ID='BridgeWRT'/" /etc/openwrt_release 2>/dev/null
sed -i "s/DISTRIB_DESCRIPTION=.*/DISTRIB_DESCRIPTION='BridgeWRT ${BRIDGEWRT_VERSION}'/" /etc/openwrt_release 2>/dev/null

# --- Hostname + Timezone ---
uci set system.@system[0].hostname='bridge-box'
uci set system.@system[0].timezone='MSK-3'
uci set system.@system[0].zonename='Europe/Moscow'
uci commit system

# --- Сеть ---
# LAN на 192.168.77.1/24 — не конфликтует с типичными роутерами
uci set network.lan.ipaddr='192.168.77.1'
uci set network.lan.netmask='255.255.255.0'
uci commit network

# --- DHCP ---
uci set dhcp.lan.start='100'
uci set dhcp.lan.limit='150'
uci set dhcp.lan.leasetime='12h'
uci commit dhcp

# --- Состояние ---
safe_write /etc/bridgebox/state 'setup'
safe_write /etc/bridgebox/wifi-mode 'down'

# --- dnsmasq оставляем включённым ---
# Нужен для: DHCP на LAN (setup), DNS hijack в AP mode, резолв hostname
# В bridge mode setup-bridge.sh НЕ отключает его (нужен для Wi-Fi AP)

# --- Автозапуск сервисов ---
/etc/init.d/bridgebox-wifi enable 2>/dev/null || true
/etc/init.d/bridgebox-agent enable 2>/dev/null || true
/etc/init.d/bridgebox-socks enable 2>/dev/null || true

exit 0
