# imagebuilder/files/etc/uci-defaults/10-bridgebox-system
#!/bin/sh
# 10-bridgebox-system — Базовая настройка системы при первой загрузке
# Выполняется один раз (uci-defaults).

BRIDGEWRT_VERSION="0.1.0"

# --- Общие утилиты ---
. /usr/lib/bridgebox/lib-common.sh

# --- BOX_ID (генерация из MAC eth0 при первой загрузке) ---
BOX_ID_FILE="/etc/bridgebox/box-id"
EXISTING_ID=$(cat "$BOX_ID_FILE" 2>/dev/null)
if [ -z "$EXISTING_ID" ] || [ "$EXISTING_ID" = "TEMPLATE" ]; then
    # Ждём eth0 (до 10 сек — для rockchip/R3S драйвер грузится асинхронно)
    MAC=""
    attempts=0
    while [ $attempts -lt 10 ]; do
        MAC=$(cat /sys/class/net/eth0/address 2>/dev/null)
        [ -n "$MAC" ] && [ "$MAC" != "00:00:00:00:00:00" ] && break
        MAC=""
        attempts=$((attempts + 1))
        sleep 1
    done
    if [ -n "$MAC" ]; then
        BOX_ID="BB-$(echo -n "$MAC" | sha256sum | cut -c1-6 | tr 'a-f' 'A-F')"
        mkdir -p /etc/bridgebox
        safe_write "$BOX_ID_FILE" "$BOX_ID"
    fi
fi

# --- Branding ---
cat > /etc/banner <<'BANNER'
░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░
░       ░░░       ░░░        ░░       ░░░░      ░░░        ░░  ░░░░  ░░       ░░░        ░
▒  ▒▒▒▒  ▒▒  ▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒  ▒▒  ▒▒▒▒▒▒▒▒  ▒▒▒▒▒▒▒▒  ▒  ▒  ▒▒  ▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒
▓       ▓▓▓       ▓▓▓▓▓▓  ▓▓▓▓▓  ▓▓▓▓  ▓▓  ▓▓▓   ▓▓      ▓▓▓▓        ▓▓       ▓▓▓▓▓▓  ▓▓▓▓
█  ████  ██  ███  ██████  █████  ████  ██  ████  ██  ████████   ██   ██  ███  ██████  ████
█       ███  ████  ██        ██       ████      ███        ██  ████  ██  ████  █████  ████
██████████████████████████████████████████████████████████████████████████████████████████

 ---------------------------------------------------------
 Transparent L2 Bridge + Remote Management
 ---------------------------------------------------------
BANNER

BOX_ID=$(cat /etc/bridgebox/box-id 2>/dev/null)
if [ -n "$BOX_ID" ] && [ "$BOX_ID" != "TEMPLATE" ]; then
    echo " Box: $BOX_ID  |  v${BRIDGEWRT_VERSION}" >> /etc/banner
else
    echo " v${BRIDGEWRT_VERSION}" >> /etc/banner
fi
echo "" >> /etc/banner

sed -i "s/DISTRIB_ID=.*/DISTRIB_ID='BridgeWRT'/" /etc/openwrt_release 2>/dev/null
sed -i "s/DISTRIB_DESCRIPTION=.*/DISTRIB_DESCRIPTION='BridgeWRT ${BRIDGEWRT_VERSION}'/" /etc/openwrt_release 2>/dev/null

# --- Hostname + Timezone (FirstBoot only — оператор может менять через mesh) ---
first_boot_only uci set system.@system[0].hostname='bridge-box'
first_boot_only uci set system.@system[0].timezone='MSK-3'
first_boot_only uci set system.@system[0].zonename='Europe/Moscow'
if is_first_boot; then
    uci commit system
fi

# --- Сеть ---
# Конфиг подложен в /etc/config/network (files overlay):
#   eth1 = LAN management (192.168.77.1), standalone, НЕ в bridge
#   eth0 = WAN, proto none
# Bridge создаётся позже через setup-bridge.sh после provisioning.
# Так оператор ВСЕГДА имеет доступ через LAN-кабель.

# --- DHCP ---
# Конфиг подложен в /etc/config/dhcp (files overlay)

# --- LED-индикация (FirstBoot only) ---
# Exhaustive board_name case — unknown board = skip с warning
if is_first_boot; then
    BOARD_NAME=$(cat /tmp/sysinfo/board_name 2>/dev/null)
    case "$BOARD_NAME" in
        friendlyarm,nanopi-r2s)
            LED_LAN="green:lan"
            LED_WAN_OK="green:wan"
            LED_WAN_ERR="red:wan"
            ;;
        friendlyarm,nanopi-r3s)
            LED_LAN="green:lan"
            LED_WAN_OK="green:wan"
            LED_WAN_ERR="red:wan"
            ;;
        *)
            logger -t bridgebox -p user.warn "Unknown board: $BOARD_NAME, LED config skipped"
            LED_LAN=""
            LED_WAN_OK=""
            LED_WAN_ERR=""
            ;;
    esac

    # WAN LED = mesh status (tailscale0 link — горит когда mesh подключен)
    if [ -n "$LED_WAN_OK" ]; then
        uci set system.led_wan=led
        uci set system.led_wan.name='mesh'
        uci set system.led_wan.sysfs="$LED_WAN_OK"
        uci set system.led_wan.trigger='netdev'
        uci set system.led_wan.dev='tailscale0'
        uci set system.led_wan.mode='link'
    fi

    # LAN LED = bridge activity (выключен в setup, setup-bridge.sh привяжет к br0)
    if [ -n "$LED_LAN" ]; then
        uci set system.led_lan=led
        uci set system.led_lan.name='bridge'
        uci set system.led_lan.sysfs="$LED_LAN"
        uci set system.led_lan.trigger='none'
        uci set system.led_lan.default='0'
    fi

    uci commit system
fi

# --- nftables ---
# nfqdns использует AF_PACKET — br_netfilter не нужен.
# fw4 убран из пакетов. Все nft правила загружаются через bridgebox-nftables:
# - QUIC drop (ruleset-post/100-bridgebox-quic-drop.nft)
# - Flow offload с задержкой 30 пакетов (ruleset-post/050-bridgebox-flow-offload.nft)
# - Overlay правила (nfqdns и т.д.)

# --- Начальное состояние (FirstBoot only) ---
if is_first_boot; then
    safe_write /etc/bridgebox/state 'setup'
    safe_write /etc/bridgebox/wifi-mode 'down'
    # Сброс boot counter — новая прошивка = чистый старт.
    # boot-state.json на overlay переживает DD-перепрошивку SD,
    # поэтому safe mode со старой прошивки может остаться.
    safe_write /etc/bridgebox/boot-state.json '{"Normal":{"failure_count":"First"}}'
fi

# --- uhttpd: слушать на всех интерфейсах в setup ---
# В setup state: 0.0.0.0:80 — доступен и через eth1, и через wlan0
# setup-bridge.sh переключит на management IP после provisioning
first_boot_only uci set uhttpd.main.listen_http='0.0.0.0:80'
if is_first_boot; then
    uci commit uhttpd
fi

# --- dnsmasq оставляем включённым ---
# Нужен для: DHCP на LAN (setup), DNS hijack в AP mode, резолв hostname
# В bridge mode setup-bridge.sh НЕ отключает его (нужен для Wi-Fi AP)

# --- Автозапуск сервисов ---
/etc/init.d/bridgebox-nftables enable 2>/dev/null || true
/etc/init.d/bridgebox-wifi enable 2>/dev/null || true
/etc/init.d/bridgebox-agent enable 2>/dev/null || true
/etc/init.d/bridgebox-socks enable 2>/dev/null || true

# --- Финализация ---
# mark_initialized ПОСЛЕДНИМ — если скрипт упал раньше, при перезагрузке повторит
if is_first_boot; then
    mark_initialized
    logger -t bridgebox "First boot initialization complete"
fi

exit 0
