# imagebuilder/files/etc/uci-defaults/10-bridgebox-system
#!/bin/sh
# 10-bridgebox-system — Базовая настройка системы при первой загрузке
# Выполняется один раз (uci-defaults).

BRIDGEWRT_VERSION="0.1.0"

# --- Branding ---
cat > /etc/banner <<'BANNER'
░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░░
░       ░░░       ░░░        ░░       ░░░░      ░░░        ░░  ░░░░  ░░       ░░░        ░
▒  ▒▒▒▒  ▒▒  ▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒  ▒▒  ▒▒▒▒▒▒▒▒  ▒▒▒▒▒▒▒▒  ▒  ▒  ▒▒  ▒▒▒▒  ▒▒▒▒▒  ▒▒▒▒
▓       ▓▓▓       ▓▓▓▓▓▓  ▓▓▓▓▓  ▓▓▓▓  ▓▓  ▓▓▓   ▓▓      ▓▓▓▓        ▓▓       ▓▓▓▓▓▓  ▓▓▓▓
█  ████  ██  ███  ██████  █████  ████  ██  ████  ██  ████████   ██   ██  ███  ██████  ████
█       ███  ████  ██        ██       ████      ███        ██  ████  ██  ████  █████  ████
██████████████████████████████████████████████████████████████████████████████████████████

 ---------------------------------------------------------
 Transparent L2 Bridge + Remote Management
 ---------------------------------------------------------
BANNER

BOX_ID=$(cat /etc/bridgebox/box-id 2>/dev/null)
if [ -n "$BOX_ID" ] && [ "$BOX_ID" != "TEMPLATE" ]; then
    echo " Box: $BOX_ID  |  v${BRIDGEWRT_VERSION}" >> /etc/banner
else
    echo " v${BRIDGEWRT_VERSION}" >> /etc/banner
fi
echo "" >> /etc/banner

sed -i "s/DISTRIB_ID=.*/DISTRIB_ID='BridgeWRT'/" /etc/openwrt_release 2>/dev/null
sed -i "s/DISTRIB_DESCRIPTION=.*/DISTRIB_DESCRIPTION='BridgeWRT ${BRIDGEWRT_VERSION}'/" /etc/openwrt_release 2>/dev/null

# --- Hostname + Timezone ---
uci set system.@system[0].hostname='bridge-box'
uci set system.@system[0].timezone='MSK-3'
uci set system.@system[0].zonename='Europe/Moscow'
uci commit system

# --- Сеть ---
# LAN на 192.168.77.1/24 — не конфликтует с типичными роутерами
uci set network.lan.ipaddr='192.168.77.1'
uci set network.lan.netmask='255.255.255.0'
uci commit network

# --- DHCP ---
uci set dhcp.lan.start='100'
uci set dhcp.lan.limit='150'
uci set dhcp.lan.leasetime='12h'
uci commit dhcp

# --- Состояние ---
echo 'setup' > /etc/bridgebox/state

# --- Автозапуск агента ---
/etc/init.d/bridgebox-agent enable 2>/dev/null || true

exit 0
