# imagebuilder/files/etc/uci-defaults/20-bridgebox-uhttpd
#!/bin/sh
# 20-bridgebox-uhttpd — Настройка uhttpd для status page
# Выполняется один раз (uci-defaults).

# Слушаем только на LAN (192.168.77.1), не на всех интерфейсах
uci set uhttpd.main.listen_http='192.168.77.1:80'
uci delete uhttpd.main.listen_https 2>/dev/null
uci set uhttpd.main.home='/www'
uci set uhttpd.main.cgi_prefix='/cgi-bin'
uci commit uhttpd

# Включаем uhttpd
/etc/init.d/uhttpd enable

exit 0
