# imagebuilder/files/etc/uci-defaults/20-bridgebox-uhttpd
#!/bin/sh
# 20-bridgebox-uhttpd — Настройка uhttpd для status page + captive portal
# Выполняется один раз (uci-defaults).

# Слушаем на 0.0.0.0 — в AP mode юзер подключается к 192.168.77.1,
# в STA mode статус-страница доступна по IP коробки
uci set uhttpd.main.listen_http='0.0.0.0:80'
uci delete uhttpd.main.listen_https 2>/dev/null
uci set uhttpd.main.home='/www'
uci set uhttpd.main.cgi_prefix='/cgi-bin'
uci commit uhttpd

# Включаем uhttpd
/etc/init.d/uhttpd enable

exit 0
