# imagebuilder/files/etc/uci-defaults/20-bridgebox-network
#!/bin/sh
# 20-bridgebox-network — Safety net для bridge конфигурации
#
# Выполняется один раз (uci-defaults).
# Проверяет что network config корректен для bridge mode.
# files/etc/config/network уже содержит правильный конфиг,
# но этот скрипт — страховка на случай если Image Builder
# перезапишет его дефолтной конфигурацией.

# Удаляем дефолтные интерфейсы OpenWrt если они есть
for iface in wan wan6 lan; do
    uci -q get network.$iface >/dev/null 2>&1 && uci delete network.$iface
done

# Удаляем дефолтный device br-lan если есть
uci -q get network.@device[0] >/dev/null 2>&1 && {
    dev_name=$(uci -q get network.@device[0].name)
    if [ "$dev_name" = "br-lan" ]; then
        uci delete network.@device[0]
    fi
}

# Проверяем наличие br0 — если нет, создаём
if ! uci -q get network.br0 >/dev/null 2>&1; then
    # Ищем device с name=br0
    has_br0=0
    idx=0
    while uci -q get network.@device[$idx] >/dev/null 2>&1; do
        name=$(uci -q get network.@device[$idx].name)
        if [ "$name" = "br0" ]; then
            has_br0=1
            break
        fi
        idx=$((idx + 1))
    done

    if [ "$has_br0" = "0" ]; then
        uci set network.br0=device
        uci set network.br0.name='br0'
        uci set network.br0.type='bridge'
        uci set network.br0.stp='0'
        uci add_list network.br0.ports='eth0'
        uci add_list network.br0.ports='eth1'
    fi
fi

# Проверяем наличие interface bridge
if ! uci -q get network.bridge >/dev/null 2>&1; then
    uci set network.bridge=interface
    uci set network.bridge.device='br0'
    uci set network.bridge.proto='dhcp'
fi

uci commit network

exit 0
