# imagebuilder/files/etc/uci-defaults/30-bridgebox-tailscale
#!/bin/sh
# 30-bridgebox-tailscale — Настройка Tailscale при первой загрузке
#
# Выполняется один раз (uci-defaults).
# Читает BOX_ID и (опционально) auth-key для авторегистрации в Headscale.

BOX_ID=$(cat /etc/bridgebox/box-id 2>/dev/null)
AUTH_KEY_FILE="/etc/bridgebox/auth-key"

# Включаем tailscale сервис
if [ -f /etc/init.d/tailscale ]; then
    /etc/init.d/tailscale enable
fi

# Hostname для Tailscale = BOX_ID в нижнем регистре
if [ -n "$BOX_ID" ] && [ "$BOX_ID" != "TEMPLATE" ]; then
    HOSTNAME=$(echo "$BOX_ID" | tr 'A-Z' 'a-z')
else
    HOSTNAME="bridgebox"
fi

# Если есть auth-key — запланировать авторегистрацию
# Авторегистрация произойдёт при первом запуске tailscale
if [ -f "$AUTH_KEY_FILE" ]; then
    AUTH_KEY=$(cat "$AUTH_KEY_FILE")

    if [ -n "$AUTH_KEY" ]; then
        # Создаём скрипт авторегистрации, который запустится после tailscale
        cat > /etc/bridgebox/tailscale-register.sh <<REGEOF
#!/bin/sh
# Одноразовый скрипт авторегистрации в Headscale
# Запускается из rc.local после первой загрузки

logger -t bridgebox-ts "Авторегистрация в Headscale как $HOSTNAME"

# Ждём запуска tailscale daemon
waited=0
while [ \$waited -lt 30 ]; do
    tailscale status >/dev/null 2>&1 && break
    sleep 1
    waited=\$((waited + 1))
done

# Регистрируемся
tailscale up \\
    --authkey="$AUTH_KEY" \\
    --hostname="$HOSTNAME" \\
    --accept-routes=false \\
    --advertise-routes=""

# Ждём IP
waited=0
while [ \$waited -lt 30 ]; do
    TS_IP=\$(tailscale ip -4 2>/dev/null)
    [ -n "\$TS_IP" ] && break
    sleep 1
    waited=\$((waited + 1))
done

if [ -n "\$TS_IP" ]; then
    logger -t bridgebox-ts "Зарегистрирован: \$TS_IP (\$HOSTNAME)"
else
    logger -t bridgebox-ts "WARN: IP не получен за 30с"
fi

# Удаляем auth-key (одноразовый)
rm -f "$AUTH_KEY_FILE"

# Удаляем себя
rm -f /etc/bridgebox/tailscale-register.sh

# Убираем из rc.local
sed -i '/tailscale-register/d' /etc/rc.local
REGEOF
        chmod +x /etc/bridgebox/tailscale-register.sh

        # Добавляем в rc.local (выполнится после всех init.d)
        if ! grep -q "tailscale-register" /etc/rc.local 2>/dev/null; then
            sed -i '/^exit 0$/d' /etc/rc.local 2>/dev/null
            echo "/etc/bridgebox/tailscale-register.sh &" >> /etc/rc.local
            echo "exit 0" >> /etc/rc.local
        fi
    fi
fi

exit 0
