# imagebuilder/test/qemu/qemu-overlay/etc/uci-defaults/99-qemu-test
#!/bin/sh
# 99-qemu-test — Настройки для QEMU тестового окружения
#
# QEMU user-mode networking (SLIRP) ожидает гостя в сети 10.0.2.0/24.
# Наш uci-defaults ставит статический 192.168.77.1 — SLIRP не может
# маршрутизировать к нему. Добавляем DHCP-клиент на wan для связности.
#
# В armsr/armv8 только один eth0. Он уже в br-lan (static 192.168.77.1).
# Создаём wan-интерфейс на eth0 (вне bridge) с DHCP — QEMU даст 10.0.2.15.
# Это позволит SSH-портфорварду работать.

# Убираем eth0 из br-lan (в QEMU нет второго порта, bridge бессмыслен)
uci del_list network.@device[0].ports='eth0' 2>/dev/null || true

# Создаём wan-интерфейс для связности с хостом
uci set network.wan=interface
uci set network.wan.device='eth0'
uci set network.wan.proto='dhcp'
uci commit network

# Dropbear: слушать на всех интерфейсах
uci set dropbear.@dropbear[0].Interface=''
uci commit dropbear

# fw4 убран из пакетов — firewall не блокирует SSH

exit 0
