# firmware/init.d/bridgebox-wifi-client
#!/bin/sh /etc/rc.common
# bridgebox-wifi-client — Автоподключение USB Wi-Fi к домашней сети
#
# [UNSTABLE] Экспериментальный. Проверен только с RTL8188EUS.
#
# Что делает:
#   1. Ждёт появления Wi-Fi phy
#   2. Создаёт wlan0, подключается через wpa_supplicant
#   3. Получает IP по DHCP
#
# Установка:
#   cp bridgebox-wifi-client /etc/init.d/
#   chmod +x /etc/init.d/bridgebox-wifi-client
#   /etc/init.d/bridgebox-wifi-client enable
#
# Конфигурация:
#   Создайте /etc/bridgebox-wpa.conf (setup-wifi-client.sh делает это автоматически)

START=98
USE_PROCD=1

WLAN="wlan0"
WPA_CONF="/etc/bridgebox-wpa.conf"
PHY_WAIT=30

boot() {
    start "$@"
}

start_service() {
    logger -t bridgebox-wifi "Запуск Wi-Fi client [UNSTABLE]"

    if [ ! -f "$WPA_CONF" ]; then
        logger -t bridgebox-wifi "WARN: $WPA_CONF не найден, пропускаем"
        return 0
    fi

    # Перезагружаем драйвер — при загрузке firmware может не успеть подгрузиться
    rmmod rtl8xxxu 2>/dev/null
    sleep 2
    modprobe rtl8xxxu

    # Ждём появления USB Wi-Fi phy (драйвер rtl8xxxu)
    PHY=""
    waited=0
    while [ $waited -lt $PHY_WAIT ]; do
        for p in /sys/class/ieee80211/phy*; do
            if [ -d "$p" ]; then
                driver=$(readlink "$p/device/driver" 2>/dev/null | xargs basename 2>/dev/null)
                if [ "$driver" = "rtl8xxxu" ]; then
                    PHY=$(basename "$p")
                    break
                fi
            fi
        done
        [ -n "$PHY" ] && break
        sleep 1
        waited=$((waited + 1))
    done

    if [ -z "$PHY" ]; then
        logger -t bridgebox-wifi "WARN: rtl8xxxu phy не появился за ${PHY_WAIT}с"
        return 0
    fi

    # Даём устройству полностью инициализироваться после modprobe
    sleep 10

    logger -t bridgebox-wifi "Найден $PHY за ${waited}с"

    # Создаём интерфейс
    iw dev "$WLAN" del 2>/dev/null || true
    iw phy "$PHY" interface add "$WLAN" type managed
    ip link set "$WLAN" up

    # Подключаемся
    killall wpa_supplicant 2>/dev/null || true
    sleep 1
    wpa_supplicant -B -i "$WLAN" -c "$WPA_CONF"
    sleep 5

    if ! iw dev "$WLAN" link | grep -q "Connected"; then
        logger -t bridgebox-wifi "WARN: не удалось подключиться к Wi-Fi"
        return 0
    fi

    # DHCP
    udhcpc -i "$WLAN" -q 2>/dev/null

    WLAN_IP=$(ip -4 addr show "$WLAN" | grep -oE 'inet [0-9.]+' | awk '{print $2}')
    logger -t bridgebox-wifi "Wi-Fi подключён: $WLAN_IP на $PHY"
}

stop_service() {
    logger -t bridgebox-wifi "Остановка Wi-Fi client"
    killall wpa_supplicant 2>/dev/null || true
    iw dev "$WLAN" del 2>/dev/null || true
    logger -t bridgebox-wifi "Wi-Fi остановлен"
}
