# imagebuilder/files/etc/init.d/bridgebox-wifi
#!/bin/sh /etc/rc.common
# bridgebox-wifi — Wi-Fi state machine
#
# При загрузке:
#   - Есть wpa.conf → подключаемся (STA mode)
#   - Нет wpa.conf  → поднимаем точку доступа (AP mode)
#
# Fallback: если STA не удался → AP mode

START=90
STOP=10

WPA_CONF="/etc/bridgebox/wpa.conf"
LOCK="/tmp/bridgebox-wifi-init.lock"

# Management интерфейс: wlan0 на реальном железе, eth2 в sandbox
MGMT_IFACE_FILE="/etc/bridgebox/mgmt-iface"

boot() {
    start
}

start() {
    # Защита от одновременного запуска (hotplug может тригернуть параллельно)
    if [ -f "$LOCK" ]; then
        logger -t bridgebox-wifi "Уже запущен (lock exists), пропускаем"
        return 0
    fi
    echo $$ > "$LOCK"

    # Ждём появления USB Wi-Fi адаптера (rtl8xxxu грузится асинхронно)
    local attempt=0
    while [ $attempt -lt 15 ]; do
        [ -n "$(ls /sys/class/ieee80211/ 2>/dev/null)" ] && break
        attempt=$((attempt + 1))
        sleep 2
    done

    if [ -z "$(ls /sys/class/ieee80211/ 2>/dev/null)" ]; then
        logger -t bridgebox-wifi "Wi-Fi адаптер не найден за 30 сек"

        if [ -f "$MGMT_IFACE_FILE" ]; then
            # Sandbox: management через ethernet (eth2)
            local mgmt_iface
            mgmt_iface=$(cat "$MGMT_IFACE_FILE")
            logger -t bridgebox-wifi "Sandbox mode: management через $mgmt_iface (DHCP)"
            # DHCP на mgmt интерфейсе уже настроен через uci (99-sandbox)
            # Ждём получения IP
            local ip_attempt=0
            local mgmt_ip=""
            while [ $ip_attempt -lt 15 ]; do
                mgmt_ip=$(ip -4 addr show "$mgmt_iface" 2>/dev/null | grep -o 'inet [0-9.]*' | cut -d' ' -f2)
                [ -n "$mgmt_ip" ] && break
                ip_attempt=$((ip_attempt + 1))
                sleep 2
            done
            if [ -n "$mgmt_ip" ]; then
                logger -t bridgebox-wifi "Sandbox: $mgmt_iface IP=$mgmt_ip"
                . /usr/lib/bridgebox/lib-common.sh
                safe_write /etc/bridgebox/wifi-mode "sta"
                rm -f "$LOCK"
                return 0
            else
                logger -t bridgebox-wifi "Sandbox: $mgmt_iface без IP за 30 сек"
                . /usr/lib/bridgebox/lib-common.sh
                safe_write /etc/bridgebox/wifi-mode "down"
                rm -f "$LOCK"
                return 1
            fi
        fi

        # Реальное железо: аварийный management на eth1
        logger -t bridgebox-wifi "Аварийный доступ через eth1 (192.168.77.1)"
        ip addr add 192.168.77.1/24 dev eth1 2>/dev/null || true
        ip link set eth1 up 2>/dev/null || true
        . /usr/lib/bridgebox/lib-common.sh
        safe_write /etc/bridgebox/wifi-mode "down"
        rm -f "$LOCK"
        return 1
    fi

    # Даём драйверу стабилизироваться
    sleep 3

    if [ -f "$WPA_CONF" ]; then
        /usr/lib/bridgebox/wifi-switch.sh restore
    else
        /usr/lib/bridgebox/wifi-switch.sh ap
    fi

    rm -f "$LOCK"
}

stop() {
    killall hostapd 2>/dev/null || true
    killall wpa_supplicant 2>/dev/null || true
    if [ -f /tmp/dnsmasq-bridgebox.pid ]; then
        kill "$(cat /tmp/dnsmasq-bridgebox.pid)" 2>/dev/null || true
    fi
    ip link set wlan0 down 2>/dev/null || true
    iw dev wlan0 del 2>/dev/null || true
}
