# imagebuilder/files/etc/init.d/bridgebox-wifi
#!/bin/sh /etc/rc.common
# bridgebox-wifi — Wi-Fi state machine
#
# При загрузке:
#   - Есть wpa.conf → подключаемся (STA mode)
#   - Нет wpa.conf  → management через eth0 (DHCP от роутера юзера)
#
# Wi-Fi credentials вводит юзер через CGI setup page.

START=90
STOP=10

WPA_CONF="/etc/bridgebox/wpa.conf"
LOCK="/tmp/bridgebox-wifi-init.lock"

# Management интерфейс: wlan0 на реальном железе, eth2 в sandbox
MGMT_IFACE_FILE="/etc/bridgebox/mgmt-iface"

boot() {
    start
}

start() {
    # Защита от одновременного запуска (hotplug может тригернуть параллельно)
    if [ -f "$LOCK" ]; then
        logger -t bridgebox-wifi "Уже запущен (lock exists), пропускаем"
        return 0
    fi
    echo $$ > "$LOCK"

    # LED: мигаем WAN LED пока ищем Wi-Fi (оператор видит что коробка загружается)
    led_blink() {
        echo "timer" > /sys/class/leds/green:wan/trigger 2>/dev/null
        echo 500 > /sys/class/leds/green:wan/delay_on 2>/dev/null
        echo 500 > /sys/class/leds/green:wan/delay_off 2>/dev/null
    }
    led_on() {
        echo "default-on" > /sys/class/leds/green:wan/trigger 2>/dev/null
    }
    led_off() {
        echo "none" > /sys/class/leds/green:wan/trigger 2>/dev/null
        echo 0 > /sys/class/leds/green:wan/brightness 2>/dev/null
    }

    led_blink
    logger -t bridgebox-wifi "Ожидание Wi-Fi адаптера..."

    # Ждём появления USB Wi-Fi адаптера (rtl8xxxu грузится асинхронно)
    local attempt=0
    while [ $attempt -lt 15 ]; do
        [ -n "$(ls /sys/class/ieee80211/ 2>/dev/null)" ] && break
        attempt=$((attempt + 1))
        sleep 2
    done

    if [ -z "$(ls /sys/class/ieee80211/ 2>/dev/null)" ]; then
        logger -t bridgebox-wifi "Wi-Fi адаптер не найден за 30 сек"

        if [ -f "$MGMT_IFACE_FILE" ]; then
            # Sandbox: management через ethernet (eth2)
            local mgmt_iface
            mgmt_iface=$(cat "$MGMT_IFACE_FILE")
            logger -t bridgebox-wifi "Sandbox mode: management через $mgmt_iface (DHCP)"
            local ip_attempt=0
            local mgmt_ip=""
            while [ $ip_attempt -lt 15 ]; do
                mgmt_ip=$(ip -4 addr show "$mgmt_iface" 2>/dev/null | grep -o 'inet [0-9.]*' | cut -d' ' -f2)
                [ -n "$mgmt_ip" ] && break
                ip_attempt=$((ip_attempt + 1))
                sleep 2
            done
            if [ -n "$mgmt_ip" ]; then
                logger -t bridgebox-wifi "Sandbox: $mgmt_iface IP=$mgmt_ip"
                . /usr/lib/bridgebox/lib-common.sh
                safe_write /etc/bridgebox/wifi-mode "sta"
                led_on
                rm -f "$LOCK"
                return 0
            else
                logger -t bridgebox-wifi "Sandbox: $mgmt_iface без IP за 30 сек"
                . /usr/lib/bridgebox/lib-common.sh
                safe_write /etc/bridgebox/wifi-mode "down"
                rm -f "$LOCK"
                return 1
            fi
        fi

        # Реальное железо без Wi-Fi: management через eth0 (DHCP от роутера юзера)
        logger -t bridgebox-wifi "Wi-Fi недоступен — management через eth0"
        . /usr/lib/bridgebox/lib-common.sh
        safe_write /etc/bridgebox/wifi-mode "down"
        led_off
        rm -f "$LOCK"
        return 0
    fi

    # Wi-Fi найден
    logger -t bridgebox-wifi "Wi-Fi адаптер найден"
    led_on

    # Даём драйверу стабилизироваться
    sleep 3

    if [ -f "$WPA_CONF" ]; then
        /usr/lib/bridgebox/wifi-switch.sh restore
    else
        # Нет wpa.conf — юзер ещё не ввёл credentials через CGI
        # Management через eth0, Wi-Fi будет настроен позже
        logger -t bridgebox-wifi "Нет wpa.conf — ожидаем настройку через CGI, management через eth0"
        . /usr/lib/bridgebox/lib-common.sh
        safe_write /etc/bridgebox/wifi-mode "down"
    fi

    rm -f "$LOCK"
}

stop() {
    killall wpa_supplicant 2>/dev/null || true
    if [ -f /tmp/udhcpc-wlan0.pid ]; then
        kill "$(cat /tmp/udhcpc-wlan0.pid)" 2>/dev/null || true
    fi
    ip link set wlan0 down 2>/dev/null || true
    iw dev wlan0 del 2>/dev/null || true
}
