# imagebuilder/files/etc/init.d/bridgebox-safe-mode
#!/bin/sh /etc/rc.common
# bridgebox-safe-mode — счётчик загрузок, переход в safe mode после 3 неудач
#
# Запускается ПЕРВЫМ (START=10), до агента и Wi-Fi.
# Считает попытки загрузки: если 3 раза подряд не было успешного heartbeat —
# переключается в safe mode (агент стартует с --safe-mode).
#
# Успешная загрузка сбрасывает счётчик (делается из bb-agent при heartbeat OK).
#
# Файлы:
#   /etc/bridgebox/boot-state.json  — персистентный счётчик (survives reboot)
#   /tmp/bridgebox-safe-mode        — флаг safe mode (volatile, читается агентом)

START=10

BOOT_STATE="/etc/bridgebox/boot-state.json"
SAFE_MODE_FLAG="/tmp/bridgebox-safe-mode"

start() {
    mkdir -p /etc/bridgebox

    # Первый boot — agent ещё не писал файл
    if [ ! -f "$BOOT_STATE" ]; then
        return 0
    fi

    local state
    state=$(cat "$BOOT_STATE" 2>/dev/null || echo '{}')

    # Формат от bb-agent (serde):
    #   {"Normal":{"failure_count":"First"}}   — нормальный режим
    #   "SafeMode"                             — safe mode
    #
    # Shell ТОЛЬКО ЧИТАЕТ. Инкремент и сброс — в agent.
    # Safe mode если: уже SafeMode ИЛИ failure_count=Third (3 неудачи подряд).

    if echo "$state" | grep -q '"SafeMode"'; then
        touch "$SAFE_MODE_FLAG"
        logger -t bridgebox-safe-mode "SAFE MODE активен — агент запустится с --safe-mode"
        return 0
    fi

    local failure_count
    failure_count=$(echo "$state" | sed -n 's/.*"failure_count"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p')

    case "$failure_count" in
        First)
            logger -t bridgebox-safe-mode "попытка загрузки 1/3"
            ;;
        Second)
            logger -t bridgebox-safe-mode "попытка загрузки 2/3"
            ;;
        Third)
            touch "$SAFE_MODE_FLAG"
            logger -t bridgebox-safe-mode "ПЕРЕХОД В SAFE MODE — 3 неудачных загрузки подряд"
            ;;
        *)
            logger -t bridgebox-safe-mode "неизвестный failure_count: $failure_count — пропуск"
            ;;
    esac
}
