# imagebuilder/files/etc/init.d/bridgebox-management
#!/bin/sh /etc/rc.common
# bridgebox-management — Автозапуск management-слоя при загрузке
#
# Что делает:
#   1. Ждёт поднятия tailscale0
#   2. Настраивает nftables: SSH только через tailscale0

START=99
USE_PROCD=1

TAILSCALE_WAIT=60

boot() {
    start "$@"
}

start_service() {
    logger -t bridgebox-mgmt "Запуск management-слоя"

    # Ждём поднятия tailscale0
    waited=0
    while [ $waited -lt $TAILSCALE_WAIT ]; do
        if ip link show tailscale0 >/dev/null 2>&1; then
            logger -t bridgebox-mgmt "tailscale0 поднят за ${waited}с"
            break
        fi
        sleep 1
        waited=$((waited + 1))
    done

    if ! ip link show tailscale0 >/dev/null 2>&1; then
        logger -t bridgebox-mgmt "WARN: tailscale0 не поднялся за ${TAILSCALE_WAIT}с, открываем SSH fallback"
        apply_firewall_fallback
        return 0
    fi

    # nftables: SSH только через tailscale0
    apply_firewall

    # Запускаем watchdog — если tailscale отвалится, откроем SSH на Wi-Fi
    procd_open_instance tailscale-watchdog
    procd_set_param command /bin/sh -c '
        while true; do
            sleep 30
            if ! tailscale status >/dev/null 2>&1; then
                logger -t bridgebox-mgmt "WARN: tailscale offline, включаем SSH fallback"
                nft delete table inet bridgebox_mgmt 2>/dev/null || true
                nft -f - <<EOF
table inet bridgebox_mgmt {
    chain input {
        type filter hook input priority 0; policy accept;
        iifname "tailscale0" tcp dport 22 accept
        iifname "wlan0" tcp dport 22 accept
        tcp dport 22 drop
    }
}
EOF
                logger -t bridgebox-mgmt "SSH fallback активен на wlan0"
                # Ждём восстановления tailscale
                while ! tailscale status >/dev/null 2>&1; do
                    sleep 30
                done
                logger -t bridgebox-mgmt "tailscale восстановлен, возвращаем strict firewall"
                nft delete table inet bridgebox_mgmt 2>/dev/null || true
                nft -f - <<EOF
table inet bridgebox_mgmt {
    chain input {
        type filter hook input priority 0; policy accept;
        iifname "tailscale0" tcp dport 22 accept
        tcp dport 22 drop
    }
}
EOF
            fi
        done
    '
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    logger -t bridgebox-mgmt "Остановка management-слоя"
    # Убираем наши правила
    nft delete table inet bridgebox_mgmt 2>/dev/null || true
    logger -t bridgebox-mgmt "Firewall правила удалены"
}

apply_firewall_fallback() {
    # Tailscale не поднялся — разрешаем SSH на Wi-Fi как аварийный вход
    nft delete table inet bridgebox_mgmt 2>/dev/null || true

    nft -f - <<'NFTABLES'
table inet bridgebox_mgmt {
    chain input {
        type filter hook input priority 0; policy accept;

        # Разрешаем SSH на tailscale0 (когда поднимется)
        iifname "tailscale0" tcp dport 22 accept

        # Fallback: SSH через Wi-Fi
        iifname "wlan0" tcp dport 22 accept

        # Блокируем SSH на остальных интерфейсах
        tcp dport 22 drop
    }
}
NFTABLES

    logger -t bridgebox-mgmt "Firewall fallback: SSH через tailscale0 + wlan0"
}

apply_firewall() {
    # Удаляем старые правила (идемпотентность)
    nft delete table inet bridgebox_mgmt 2>/dev/null || true

    nft -f - <<'NFTABLES'
table inet bridgebox_mgmt {
    chain input {
        type filter hook input priority 0; policy accept;

        # Разрешаем SSH на tailscale0
        iifname "tailscale0" tcp dport 22 accept

        # Блокируем SSH на всех остальных интерфейсах
        tcp dport 22 drop
    }
}
NFTABLES

    logger -t bridgebox-mgmt "Firewall применён: SSH только через tailscale0"
}
