# imagebuilder/files/etc/init.d/bridgebox-agent
#!/bin/sh /etc/rc.common
# bridgebox-agent — регистрация устройства при загрузке
#
# Запускается ПОСЛЕ bridgebox-wifi (START=90).
# Если Wi-Fi уже в STA mode — регистрируется сразу.
# Если Wi-Fi в AP mode — регистрация произойдёт после ввода Wi-Fi credentials
# (вызывается из wifi-setup CGI).

START=99

boot() {
    start
}

start() {
    # Проверяем Wi-Fi mode
    local wifi_mode
    wifi_mode=$(cat /etc/bridgebox/wifi-mode 2>/dev/null || echo "down")

    if [ "$wifi_mode" != "sta" ]; then
        logger -t bb-agent "Management не готов ($wifi_mode), регистрация отложена"
        return 0
    fi

    # Management интерфейс
    local mgmt_iface="wlan0"
    if [ -f /etc/bridgebox/mgmt-iface ]; then
        mgmt_iface=$(cat /etc/bridgebox/mgmt-iface)
    fi

    # Ждём IP на management интерфейсе
    local attempt=0
    local ip=""
    while [ $attempt -lt 10 ]; do
        ip=$(ip -4 addr show "$mgmt_iface" 2>/dev/null | grep -o 'inet [0-9.]*' | cut -d' ' -f2)
        [ -n "$ip" ] && break
        attempt=$((attempt + 1))
        sleep 2
    done

    if [ -z "$ip" ]; then
        logger -t bb-agent "$mgmt_iface без IP — ретрай через heartbeat"
        return 1
    fi

    echo "nameserver 8.8.8.8" > /tmp/resolv.conf  # fallback DNS до DHCP

    # Регистрация
    /usr/bin/bb-agent register 2>&1 | logger -t bb-agent

    # Mesh — поднимаем сразу после регистрации (не ждём heartbeat через 5 мин)
    /usr/bin/bb-agent ensure-mesh 2>&1 | logger -t bb-agent
}
