# imagebuilder/files/etc/hotplug.d/usb/30-bridgebox-wifi
#!/bin/sh
# 30-bridgebox-wifi — Уведомляет Wi-Fi state machine о появлении USB адаптера
#
# НЕ создаёт интерфейс сам — это делает wifi-switch.sh
# Только тригерит запуск, если state machine ещё не работает

[ "$ACTION" = "add" ] || exit 0

# Если wlan0 уже есть — адаптер уже подхвачен
[ -d /sys/class/net/wlan0 ] && exit 0

# Ждём появления phy (драйвер грузится асинхронно)
attempts=0
while [ $attempts -lt 5 ]; do
    [ -n "$(ls /sys/class/ieee80211/ 2>/dev/null)" ] && break
    attempts=$((attempts + 1))
    sleep 1
done

[ -z "$(ls /sys/class/ieee80211/ 2>/dev/null)" ] && exit 0

# Тригерим state machine если она ещё не запущена
WIFI_MODE=$(cat /etc/bridgebox/wifi-mode 2>/dev/null || echo "down")
LOCK="/tmp/bridgebox-wifi-init.lock"
if [ "$WIFI_MODE" = "down" ] && [ ! -f "$LOCK" ]; then
    logger -t bridgebox-wifi "USB Wi-Fi обнаружен (hotplug), запуск state machine"
    /etc/init.d/bridgebox-wifi start &
fi
